# config/crontab
*/1 * * * * /usr/local/bin/node /usr/src/iot_raspi/dist/index.js
