# scripts/crontab
0,15,30,45  * * * * /usr/local/bin/node /usr/src/iot-raspi/index.js
