# scripts/crontab
0,15,30,45 * * * * /usr/src/iot-raspi-sensors/start-node.sh
