# scripts/crontab
0,15,30,45 * * * * /usr/src/iot_raspi/scripts/start_node.sh
