# scripts/crontab
*/1 * * * * /usr/src/iot-raspi/scripts/start_node.sh;
