# scripts/crontab
*/1 * * * * /usr/src/iot-raspi/scripts/start-node.sh
